# python-side/resources/.config/autostart/spectacle-fix.desktop
[Desktop Entry]
Type=Application
Name=Spectacle-Fix
Exec=/usr/bin/spectacle_fix